# hubcorner.service

[Unit]
Description=Hub Corner Community Platform
After=network.target

[Service]
WorkingDirectory=/opt/hubcorner
ExecStart=/opt/hubcorner/hubcorner
User=hubcorner
Group=hubcorner
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=hubcorner
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
